# django-api/crontab
*/5 * * * * /usr/local/bin/python /api/manage.py toggle_overdue_tasks
